# look_at_the_tv.service
[Service]
ExecStart=/usr/bin/look_at_the_tv
Restart=always

[Install]
WantedBy=default.target
